# mysqld.conf
[program:mysql]
process_name = mysqld_safe
directory = /var/lib/mysql
command = /usr/bin/mysqld_safe
startsecs = 0
autorestart = false
